# etc/opt/MEET/cron.daily/config-update
#!/bin/bash

#
# Update configuration files.
# [maxg]
#

( cd / && git pull origin master --quiet )
