# etc/opt/MEET/cron.daily/iceweasel
#!/bin/sh

update-alternatives --set gnome-www-browser /usr/bin/iceweasel
update-alternatives --set x-www-browser /usr/bin/iceweasel
