# etc/opt/MEET/cron.daily/meet-configs-up
#!/bin/bash

#
# Update configuration files.
# [maxg]
#

( cd / && git pull --quiet )
