# etc/cron.daily/meet
#!/bin/sh

run-parts --report /etc/opt/MEET/cron.daily
